# filament-monitor.service
[Unit]
Description=Filament motion & runout monitor (Marlin / Prusa firmware)
After=network.target

[Service]
Type=simple

# Provide a writable /run/filmon/ directory for the local control socket.
RuntimeDirectory=filmon
RuntimeDirectoryMode=0755

# Edit these paths/arguments for your installation.
# Recommended: create a venv in the project directory and reference its python.
WorkingDirectory=/opt/filament-monitor
ExecStart=/opt/filament-monitor/.venv/bin/python /opt/filament-monitor/filament-monitor.py -p /dev/ttyACM0 --control-socket /run/filmon/filmon.sock

Restart=on-failure
RestartSec=2

# Send logs to journald
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
